# doc/reference/benchmark-server-go/systemd/benchmark-server.socket
[Unit]
Description=Socket for benchmark server.

[Socket]
ListenStream=0.0.0.0:80
NoDelay=true

[Install]
WantedBy=sockets.target
