# doc/reference/benchmark-server-althea/systemd/benchmark-server-althea.service
[Unit]
Description=test server.
After=network.target
Requires=benchmark-server-althea.socket

[Service]
ExecStart=/usr/local/bin/benchmark-server-althea
#User=root
#Group=root
Restart=on-failure

[Install]
WantedBy=multi-user.target
